# src/resourced/resourced.socket
[Unit]
Description=Resourced socket

[Socket]
ListenStream=/tmp/proc_info
SocketMode=0777
SmackLabelIPIn=*
SmackLabelIPOut=@
