# packaging/resourced_swapoff.service
[Unit]
Description=Resource management daemon
After=tizen-system.target
Wants=tizen-system.target

[Service]
Type=simple
EnvironmentFile=/run/tizen-system-env
ExecStart=/usr/bin/resourced -s 1 -u 60
Restart=always
RestartSec=0

[Install]
WantedBy=multi-user.target
